# systemd/statiegeld-web.service
[Unit]
Description=Statiegeld web server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/statiegeld
Environment=ADMIN_USERNAME=admin
Environment=ADMIN_PASSWORD=changeme
Environment=API_KEY=statiegeld-scanner
Environment=SECRET_KEY=changeme-to-a-random-string
ExecStart=/home/pi/statiegeld/.venv/bin/uvicorn statiegeld.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
